# pmi-workers/cron/crontab
# pmi-workers supercronic crontab.
#
# Each line: <min> <hour> <dom> <month> <dow> <command>
# Trailing `&` keeps supercronic from blocking the schedule on a long-running
# job — same pattern as micah-job-executor/cron/crontab.
#
# Jobs are resolved via the `run-job` console script (see pyproject.toml
# [project.scripts]) which dispatches through `pmi_workers.registry`.

# Hourly: re-score every current PMI definition.
0 * * * * run-job hourly &

# Tier 3 (CORR-5.2): every 15 min, re-score only the indexes whose component
# prices drifted >= PMI_DRIFT_THRESHOLD_PCT points vs ~lookback ago. Cheap
# no-op when nothing drifted; keeps scores fresher than the hourly full pass.
*/15 * * * * run-job reeval-drifted &

# Daily 10:00 UTC: heavier maintenance (placeholder — currently same as hourly).
0 10 * * * run-job daily &

# --- OpenAI Batch API (CORR-5.3, half-price nightly recompute) --------------
# Only meaningful when factors are bound to gpt-*/openai/* models AND
# OPENAI_API_KEY is set; with an all-ollama registry both jobs no-op cleanly.
0 2 * * * run-job llm-batch-submit &
*/30 * * * * run-job llm-batch-poll &

# ─── Polymarket deep-signal cycle jobs (CORR-4.1 ~ 4.4) ─────────────────────
# These shell out to the pmi-ingest CLI. The pmi-workers image does NOT bundle
# pmi-ingest's deps (httpx/websockets/web3) — the lines below are commented as
# the intended cadence; in compose, run them as dedicated sibling services
# (see `ingest-clob` / `ingest-chain` / `ingest-ws` ideas in
# `pmi-ingest/README.md`), or uncomment after rebuilding pmi-workers with
# pmi-ingest installed.
#
#   *  *    * * *   pmi-ingest clob --once   &   # orderbook depth, every minute
#   */2 *   * * *   pmi-ingest chain --once  &   # Polygon log indexer
#   30 6,18 * * *   pmi-ingest cohort        &   # trader cohort rollup, 2× daily
#   15 *    * * *   pmi-ingest uma           &   # UMA projection, hourly
#   (ws)            pmi-ingest ws            &   # WS feed — own long-lived process
