# pmi-workers/cron/crontab
# pmi-workers supercronic crontab.
#
# Each line: <min> <hour> <dom> <month> <dow> <command>
#
# CORR-4.6 (cron → worker): schedule lines ENQUEUE onto the Postgres job queue
# (`core_jobs`); the long-lived `pmi-worker` compose service claims + executes
# (retry / heartbeat / dedupe live there). `pmi-workers enqueue` dedupes by
# job name+args, so a beat landing while the previous tick is still
# queued/running is a no-op instead of stacking — which also retires the old
# trailing-`&` pattern (enqueue is a sub-second INSERT).
#
# One-shot manual in-process execution still works: `run-job <name>`.

# Hourly: re-score every current PMI definition.
0 * * * * pmi-workers enqueue hourly

# Tier 3 (CORR-5.2): every 15 min, re-score only the indexes whose component
# prices drifted >= PMI_DRIFT_THRESHOLD_PCT points vs ~lookback ago. Cheap
# no-op when nothing drifted; keeps scores fresher than the hourly full pass.
*/15 * * * * pmi-workers enqueue reeval-drifted

# Daily 10:00 UTC: heavier maintenance (placeholder — currently same as hourly).
0 10 * * * pmi-workers enqueue daily

# --- OpenAI Batch API (CORR-5.3, half-price nightly recompute) --------------
# Only meaningful when factors are bound to gpt-*/openai/* models AND
# OPENAI_API_KEY is set; with an all-ollama registry both jobs no-op cleanly.
0 2 * * * pmi-workers enqueue llm-batch-submit
*/30 * * * * pmi-workers enqueue llm-batch-poll

# ─── Polymarket deep-signal cycle jobs (CORR-4.1 ~ 4.4) ─────────────────────
# These shell out to the pmi-ingest CLI. The pmi-workers image does NOT bundle
# pmi-ingest's deps (httpx/websockets/web3) — the lines below are commented as
# the intended cadence; in compose, run them as dedicated sibling services
# (see `ingest-clob` / `ingest-chain` / `ingest-ws` ideas in
# `pmi-ingest/README.md`), or uncomment after rebuilding pmi-workers with
# pmi-ingest installed.
#
#   *  *    * * *   pmi-ingest clob --once   &   # orderbook depth, every minute
#   */2 *   * * *   pmi-ingest chain --once  &   # Polygon log indexer
#   30 6,18 * * *   pmi-ingest cohort        &   # trader cohort rollup, 2× daily
#   15 *    * * *   pmi-ingest uma           &   # UMA projection, hourly
#   (ws)            pmi-ingest ws            &   # WS feed — own long-lived process
